# git-xopp.path
[Unit]
Description="Monitor a Xopp folder for changes"
After=graphcal.target
Requisite=graphcal.target

[Path]
PathModified="Directory"
Unit=git-xopp.service
TriggerLimitIntervalSec=10min
TriggerLimitBurst=5
