# git-xopp.service
[Unit]
Description="Run script to commit changes to Git"

[Service]
ExecStart=./git-xopp.sh Directory
